# .config/systemd/user/dotfiles-git-fetch.timer
[Unit]
Description=Fetch dotfile updates.

[Timer]
OnBootSec=15min
OnUnitInactiveSec=60min

[Install]
WantedBy=timers.target
