# .config/systemd/user/session-priority.slice
[Unit]
Description=High-priority Session Services

[Slice]
IOWeight=1000
CPUWeight=1000
