# .config/systemd/user/pipewire-pulse.service.d/slice.conf
[Service]
Slice=session-priority-pipewire.slice
