# .config/systemd/user/org.gnome.Shell@wayland.service.d/ssh-agent.conf
[Unit]
After=ssh-agent.service
Wants=ssh-agent.service
