# .config/systemd/user/gsd-color.service.d/override.conf
[Unit]
# if colord gets restarted, screen color management is broken until gsd-color
# is restarted. therefore, ensure we can restart it as needed.
RefuseManualStop=false
RefuseManualStart=false
